# infra/loki/bgs-loki.service
[Unit]
Description=Boardgamers Loki logging stack (podman-compose)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=bgs
Group=bgs
WorkingDirectory=/home/bgs/boardgamers-mono/infra/loki
ExecStart=/usr/bin/podman-compose up --no-recreate
ExecStop=/usr/bin/podman-compose down
Restart=always
RestartSec=10
TimeoutStartSec=120

[Install]
WantedBy=multi-user.target
